# src/pico_control_V0.5/scripts/vr-arm-control.service
[Unit]
Description=VR Arm Remote Control Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=aidlux
WorkingDirectory=/home/aidlux/pico/pico_control_V0.5
ExecStart=/bin/bash -c "source /opt/ros/humble/setup.bash && /usr/bin/python3 /home/aidlux/pico/pico_control_V0.5/vr_arm_control_api.py"
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=3
TimeoutStopSec=15

# 环境变量
Environment=PYTHONPATH=/home/aidlux/pico/pico_control_V0.5:/opt/ros/humble/lib/python3.10/site-packages
Environment=ROS_DOMAIN_ID=0
Environment=AMENT_PREFIX_PATH=/opt/ros/humble
Environment=CMAKE_PREFIX_PATH=/opt/ros/humble
Environment=COLCON_PREFIX_PATH=/opt/ros/humble
Environment=LD_LIBRARY_PATH=/opt/ros/humble/lib
Environment=PATH=/opt/ros/humble/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PKG_CONFIG_PATH=/opt/ros/humble/lib/pkgconfig
Environment=PYTHONPATH=/opt/ros/humble/lib/python3.10/site-packages:/home/aidlux/pico/pico_control_V0.5

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vr-arm-control

[Install]
WantedBy=multi-user.target